# configs/N1/files-23.05/files/etc/crontabs/root

5 6 * * 6 sleep 5 && touch /etc/banner && reboot
